# init_scripts/arke_upstart.conf
description "arke server monitoring agent"
author      "Dominic LoBue <dominic@geodelic.com>"


start on startup
stop on shutdown

console none
respawn

exec /usr/bin/arke -n

