# init_scripts/arke.conf
description "arke server monitoring agent"
author      "Dominic LoBue <dominic@geodelic.com>"


start on startup
stop on shutdown

kill timeout 20
console none
respawn

exec /usr/bin/arke

